# supervisord.develop.conf
[unix_http_server]
file=./tmp/supervisor.sock   ; (the path to the socket file)

[inet_http_server]         ; inet (TCP) server disabled by default
port=*:7010        ; (ip_address:port specifier, *:port for all iface)

[supervisord]
logfile=./logs/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=./tmp/pids/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:7010 ; use an http:// url to specify an inet socket

[program:mmmobile_smpp_transport]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --worker_class=vumi.workers.smpp.transport.SmppTransport
    --vhost=/develop
    --config=./config/smpp_mmmobile.yaml
    --config_smpp_offset=%(process_num)s
    --config_smpp_increment=10
environment=DJANGO_SETTINGS_MODULE=environments.develop
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true


[program:smpp_sim_worker]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --worker_class=vumi.workers.smpp.transport.SmppTransport
    --vhost=/develop
    --config=./config/smpp_sim.yaml
    --config_smpp_offset=%(process_num)s
    --config_smpp_increment=10
environment=DJANGO_SETTINGS_MODULE=environments.develop
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:vumi_test_smsc_worker]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --worker_class=vumi.workers.smpp.service.SmppService
    --vhost=/develop
    --config=./config/develop.smpp.vumitestsmsc.yaml
environment=DJANGO_SETTINGS_MODULE=environments.develop
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:vumi_test_esme_worker]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --worker_class=vumi.workers.smpp.transport.SmppTransport
    --vhost=/develop
    --config=./config/develop.smpp.vumitestsmsc.yaml
    --config_smpp_offset=%(process_num)s
    --config_smpp_increment=10
environment=DJANGO_SETTINGS_MODULE=environments.develop
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:smpp_worker]
numprocs=1
numprocs_start=2
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --worker_class=vumi.workers.smpp.transport.SmppTransport
    --vhost=/develop
    --config=./config/smpp.yaml
    --config_smpp_offset=%(process_num)s
    --config_smpp_increment=10
environment=DJANGO_SETTINGS_MODULE=environments.develop
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:debatch_worker]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --worker_class=vumi.workers.smpp.worker.SMSBatchWorker
    --vhost=/develop
environment=DJANGO_SETTINGS_MODULE=environments.develop
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:keyword_worker]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --worker_class=vumi.workers.smpp.worker.SMSKeywordWorker
    --vhost=/develop
    --config=./config/smpp.yaml
environment=DJANGO_SETTINGS_MODULE=environments.develop
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:receipt_worker]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --worker_class=vumi.workers.smpp.worker.SMSReceiptWorker
    --vhost=/develop
    --config=./config/smpp.yaml
environment=DJANGO_SETTINGS_MODULE=environments.develop
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:receipt_sim_worker]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --worker_class=vumi.workers.smpp.worker.SMSReceiptWorker
    --vhost=/develop
    --config=./config/smpp_sim.yaml
environment=DJANGO_SETTINGS_MODULE=environments.develop
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:webapp]
numprocs=1
numprocs_start=0
process_name=%(program_name)s_700%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_700%(process_num)s.pid
    vumi_webapp
    --port=700%(process_num)s
    --django-settings=environments.develop
stdout_logfile=./logs/%(program_name)s_700%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_700%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:smpp_simulator]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_700%(process_num)s
command=java 
    -Djava.net.preferIPv4Stack=true 
    -Djava.util.logging.config.file=conf/logging.properties 
    -jar smppsim.jar 
    ../../../config/smppsim.develop.props
stdout_logfile=./logs/%(program_name)s_700%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_700%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true
directory=./utils/smppsim/SMPPSim/

[program:opera_sms_transport]
numprocs=4
numprocs_start=0
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --worker_class=vumi.workers.opera.transport.OperaSMSTransport
    --vhost=/develop
    --config=./config/opera.yaml
environment=DJANGO_SETTINGS_MODULE=environments.develop
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:opera_receipt_transport]
numprocs=1
numprocs_start=0
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --worker_class=vumi.workers.opera.transport.OperaReceiptTransport
    --vhost=/develop
    --config=./config/opera.yaml
environment=DJANGO_SETTINGS_MODULE=environments.develop
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true


[program:cellulant_transport]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --worker_class=vumi.workers.httprpc.transport.HttpRpcTransport
    --vhost=/develop
    --config=./config/cellulant.yaml
environment=DJANGO_SETTINGS_MODULE=environments.develop
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:mmmobile_transport]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --worker_class=vumi.workers.httprpc.transport.HttpRpcTransport
    --vhost=/develop
    --config=./config/http_mmmobile.yaml
environment=DJANGO_SETTINGS_MODULE=environments.develop
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:integrat_inbound_transport]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --worker_class=vumi.workers.integrat.transport.IntegratInboundTransport
    --vhost=/develop
    --config=config/integrat_transport.yaml
environment=DJANGO_SETTINGS_MODULE=environments.development
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:integrat_outbound_transport]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --worker_class=vumi.workers.integrat.transport.IntegratOutboundTransport
    --vhost=/develop
    --config=config/integrat_transport.yaml
environment=DJANGO_SETTINGS_MODULE=environments.development
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:integrat_worker]
numprocs=1
numprocs_start=0
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --worker_class=vumi.workers.integrat.worker.IntegratWorker
    --config=config/integrat_worker.yaml
    --vhost=/develop
environment=DJANGO_SETTINGS_MODULE=environments.development
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true 

[program:integrat_txtalert_worker]
numprocs=1
numprocs_start=0
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --worker_class=vumi.campaigns.txtalert.USSDBookingTool
    --config=config/integrat_txtalert_worker.yaml
    --vhost=/develop
environment=DJANGO_SETTINGS_MODULE=environments.development
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true 

[program:integrat_twitter_worker]
numprocs=1
numprocs_start=0
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --worker_class=vumi.campaigns.twitter.TwitterUSSDTransport
    --config=config/twitter_ussd.yaml
    --vhost=/develop
environment=DJANGO_SETTINGS_MODULE=environments.development
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true 

[program:integrat_wikipedia_worker]
numprocs=1
numprocs_start=0
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --worker_class=vumi.campaigns.wikipedia.WikipediaWorker
    --config=config/wikipedia_ussd.yaml
    --vhost=/develop
environment=DJANGO_SETTINGS_MODULE=environments.development
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true 

[program:integrat_tictactoe_worker]
numprocs=1
numprocs_start=0
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --worker_class=vumi.workers.integrat.games.TicTacToeWorker
    --config=config/tictactoe_ussd.yaml
    --vhost=/develop
environment=DJANGO_SETTINGS_MODULE=environments.development
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true

[program:integrat_rockpaperscissors_worker]
numprocs=1
numprocs_start=0
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --worker_class=vumi.workers.integrat.games.RockPaperScissorsWorker
    --config=config/rockpaperscissors_ussd.yaml
    --vhost=/develop
environment=DJANGO_SETTINGS_MODULE=environments.development
stdout_logfile=./logs/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./logs/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
autorestart=true
